# conf/supervisor.ini
[program:l2]
command=lighttpd2 -c /srv/l2/angel.conf
autorestart=true
startretries=9999

[program:php]
command=spawn-fcgi -n -u apache -g apache -s /tmp/php.sock -S -M 0666 -C 32 -F 1 -P /tmp/php.pid -- /usr/bin/php-cgi

[program:php-dev]
command=spawn-fcgi -n -u apache -g apache -s /tmp/php-dev.sock -S -M 0666 -C 2 -P /tmp/php-dev.pid -- /usr/bin/php-cgi

[program:ws]
command=gunicorn -b 0.0.0.0:8081 -k "geventwebsocket.gunicorn.workers.GeventWebSocketWorker" wsgi-debug:application
directory=/srv/src/data.fm/py
user=nobody

[program:stunnel]
command=stunnel /srv/src/data.fm/conf/stunnel.conf
